# tools/cookie-receiver.service
[Unit]
Description=Sahibinden Cookie Receiver (cookie bridge endpoint)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /home/cagr/Masaüstü/upu-platform/tools/cookie-receiver.mjs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
